# deploy/supervisor/supervisord.conf
[inet_http_server]          ; inet (TCP) server disabled by default
port=*:9001                 ; ip_address:port specifier, *:port for all iface
username=mpedigree
password=mpedigree

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true                ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200


; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work.  Additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.

[supervisorctl]
serverurl=http://*:9001      ; use an http:// url to specify an inet socket
username=mpedigree
password=mpedigree

[program:uwsgi]
command=uwsgi --ini /Invoice_Management_System/deploy/django/invoice.ini
directory=/Invoice_Management_System
autostart=true
autorestart=true
stderr_logfile = /Invoice_Management_System/deploy/supervisor/error.log
stdout_logfile = /Invoice_Management_System/deploy/supervisor/output.log
